# link/.config/systemd/user/fcitx5.service
[Unit]

[Service]
ExecStart=/usr/bin/fcitx5 
Restart=always
