# link/.config/systemd/user/fcitx5.service
[Unit]
Wants=dbus.service

[Service]
ExecStart=/usr/bin/fcitx5 
Restart=always
